# .docker/php-fpm/supervisord.d/cron.conf
[program:cron]
process_name=%(program_name)s_%(process_num)02d
command=cron -f
autostart=true
autorestart=true
numprocs=1
; user=www-data
redirect_stderr=true
