# .docker/php-fpm/supervisord.d/php-fpm.conf
[program:php-fpm]
command=/usr/local/sbin/php-fpm
numprocs=1
autostart=true
autorestart=true
priority=100
redirect_stderr=true
